# RAWDATA/INBOUND/SCRIPTS/CSVtoTXT.service
[Unit]
Description= Servicio de Ejecución de Script CSVtoTXT

[Service]
Type=simple
ExecStart=../SCRIPTS/CSVtoTXT.py